# onlymatt-gateway/onlymatt-gateway.service
[Unit]
Description=OnlyMatt Gateway
After=network.target

[Service]
WorkingDirectory=/opt/onlymatt-gateway
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/onlymatt-gateway/.env
ExecStart=/opt/onlymatt-gateway/.venv/bin/uvicorn gateway:app --host 0.0.0.0 --port 5059 --workers 2
Restart=always
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
